# zed/service/adcserver.service
[Unit]
Description=ADC Data Server (FastAPI + iio_reader)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/root/start.sh
Restart=always
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
